# systemd/ledmatrix-stop.service
[Unit]
Description=Stop LED Matrix

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop ledmatrix.service