# systemd/temp-shutdown.service
[Unit]
Description=Check Raspberry Pi temperature

[Service]
Type=oneshot
ExecStart=/usr/local/bin/temp_shutdown.sh